# agent/systemd/blueforge-agent.service
[Unit]
Description=BlueForge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Sesuaikan path instalasi agen
WorkingDirectory=/opt/blueforge-agent
ExecStart=/usr/bin/python3 /opt/blueforge-agent/main.py
Restart=always
RestartSec=5
User=root
# Auto-restart bila crash (TDD §39 Risk Register)

[Install]
WantedBy=multi-user.target
